# 99-media-magic.rules
KERNEL=="sr0",SUBSYSTEM=="block",ACTION=="change",RUN+="/bin/systemctl start media-magic"
