# media-magic.service
[Unit]
Description=Auto CD/DVD Rip & Play

[Service]
Type=simple
ExecStart=/usr/local/bin/media-magic.sh
